# scripts/deepseek-ocr.service
[Unit]
Description=DeepSeek-OCR Web Service
After=network.target

[Service]
Type=simple
User=rongheng
WorkingDirectory=/home/rongheng/python_project/deepseek-ocr
ExecStart=/home/rongheng/python_project/deepseek-ocr/.venv/bin/deepseek-ocr serve --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=10
Environment=CUDA_VISIBLE_DEVICES=0
Environment=HF_ENDPOINT=https://hf-mirror.com

[Install]
WantedBy=multi-user.target
